# pl-bb-sd-card-falcon/epdc/bin/bluetooth/ble-autoconnect.service
[Unit]
Description=Autoconnect BLE devices from config

[Service]
WorkingDirectory=/mnt/data/epdc/bin/bluetooth/
User=debian
Environment=PATH=/home/debian/.local/lib/python3.9/site-packages:/home/debian/.local/bin
ExecStart=/usr/bin/python3 ble-autoconnect.py

[Install]
WantedBy=multi-user.target
